# dockerfiles/crontab
0 12 * * * /usr/bin/certbot renew --quiet